# packages/witticism/witticism-0.4.2.tar.gz/witticism-0.4.2/witticism.service
[Unit]
Description=Witticism Voice Transcription Service
After=graphical-session.target

[Service]
Type=simple
User=aaronontheweb
WorkingDirectory=/home/aaronontheweb/repositories/ai-playground/witticism
Environment="PYTHONPATH=/home/aaronontheweb/repositories/ai-playground/witticism/src"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/aaronontheweb/.Xauthority"
ExecStart=/home/aaronontheweb/repositories/ai-playground/witticism/venv/bin/python /home/aaronontheweb/repositories/ai-playground/witticism/src/witticism/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target